# .config/systemd/user/backup.timer
[Unit]
Description=Run backup hourly
[Timer]
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target
