# .config/systemd/user/backup.service
[Unit]
Description=backup using restic
[Service]
Environment="DISPLAY=:0"
ExecStart=%h/bin/backup.sh
Type=oneshot
